# scripts/client_node_setup/mediswarm-vpn.service
[Unit]
Description=MediSwarm VPN Connection (GoodAccess/OpenVPN)
Documentation=https://github.com/KatherLab/MediSwarm
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The config path is set by setup_vpntunnel.sh -s during installation.
# It copies the .ovpn file to /etc/openvpn/client/mediswarm.conf and
# injects keepalive + credential settings.
ExecStart=/usr/sbin/openvpn --config /etc/openvpn/client/mediswarm.conf
Restart=always
RestartSec=15
TimeoutStopSec=10

# Security hardening
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=false
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_DAC_OVERRIDE
LimitNPROC=10

[Install]
WantedBy=multi-user.target
